# sourdough.service
[Unit]
Description=Sourdough Bread Logger Server
After=network.target

[Service]
Type=simple
User=mdeckert
WorkingDirectory=/home/mdeckert/sourdough
ExecStart=/home/mdeckert/sourdough/bin/sourdough-server
Restart=always
RestartSec=10

# Environment variables
Environment="SOURDOUGH_PORT=8080"
Environment="SOURDOUGH_DATA_DIR=/home/mdeckert/sourdough/data"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sourdough

[Install]
WantedBy=multi-user.target
